# rootfs/etc/services.d/shairport-sync/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Shairport Sync
# Starts the Shairport Sync service
# ==============================================================================
# shellcheck disable=SC1091


# Wait for Avahi to become available
s6-svwait -u -t 5000 /var/run/s6/services/avahi

# Wait at least 5 seconds before starting
# Avahi might need some time.
sleep 5

declare airplay_name

if bashio::config.has_value 'airplay_name'; then
    airplay_name=$(bashio::config 'airplay_name')
else
    airplay_name=$(bashio::info.hostname)
fi

exec shairport-sync -m avahi -a "$airplay_name"
